# upload-worker/deploy/hacknet-upload.service
[Unit]
Description=Hacknet file server (uploads + downloads to local disk)
After=network.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=/usr/bin/node __INSTALL_DIR__/src/server.mjs
Restart=always
RestartSec=5
TimeoutStopSec=600

[Install]
WantedBy=multi-user.target
